# .config/supervisord.conf
[supervisord]
logfile=/srv/project/.log/supervisor.log

[program:nginx]
command=nginx -g "daemon off;"
